# build/rootfs/etc/gearbox/services/memcached/finish
#!/bin/bash
echo "memcached terminated"
